# server-configuration/upstart-forever-ilovestage.conf
# node-example.conf

description "Start ilovestage.co.uk Node servers with Upstart and Forever"

# Start up when the system hits any normal runlevel, and
# shuts down when the system goes to shutdown or reboot.
start on filesystem or runlevel [2345]
stop on runlevel [06]

# IMPORTANT: You will want to use this with Forever. It
# tells Upstart that forever is going to fork after it
# starts.
expect fork

# This monitors Forever, which seems gratuitous.
# TIP: Comment these out while debugging your script.
respawn
respawn limit 5 30

# Send error messages to the console. Useful for debugging.
console output

# ubuntu is a very low-privileged user. We run as
# this user instead of as root (the default user).
setuid ubuntu

# The user's home directory
env HOME=/home/ubuntu
env NODE_ENV=production
env PORT=5050
#env APP_PATH="/var/www/ilovestage.co.uk"
#env APP_NAME="app-api.js"

#chdir /var/www/ilovestage.co.uk

# Now we run Forever, telling it to write its logs to
# /tmp, and to run the script /var/www/ilovestage.com/app-api.js
script
  exec 2>>/var/www/ilovestage.co.uk/log/upstart.log 1>>/var/www/ilovestage.co.uk/log/upstart.log
  set -x

  exec /usr/local/bin/forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk /var/www/ilovestage.co.uk/app-api.js
  #exec /usr/local/bin/forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk /var/www/ilovestage.co.uk/app-api.js
  #exec /usr/local/bin/forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk /var/www/ilovestage.co.uk/app-api.js
  #exec /usr/local/bin/forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk /var/www/ilovestage.co.uk/app-api.js
end script

pre-stop script
  exec 2>>/var/www/ilovestage.co.uk/log/upstart.log 1>>/var/www/ilovestage.co.uk/log/upstart.log
  set -x

  exec /usr/local/bin/forever stop /var/www/ilovestage.co.uk /var/www/ilovestage.co.uk/app-api.js
end script
