# server-configuration/upstart-forever-ilovestage.conf
# node-example.conf

description "Start ilovestage.co.uk Node servers with Upstart and Forever"

# Start up when the system hits any normal runlevel, and
# shuts down when the system goes to shutdown or reboot.
start on filesystem or runlevel [2345]
stop on runlevel [06]

# IMPORTANT: You will want to use this with Forever. It
# tells Upstart that forever is going to fork after it
# starts.
expect fork

# This monitors Forever, which seems gratuitous.
# TIP: Comment these out while debugging your script.
respawn
respawn limit 5 30

# Send error messages to the console. Useful for debugging.
console output

# ubuntu is a very low-privileged user. We run as
# this user instead of as root (the default user).
setuid ubuntu

# The user's home directory
env HOME=/home/ubuntu

# Now we run Forever, telling it to write its logs to
# /tmp, and to run the script /var/www/ilovestage.com/app-api.js
script
  #exec 2>>/var/log/upstart/forever-ilovestage.log
  exec 2>>/var/www/ilovestage.co.uk/log/upstart.log
  set -x
  #export NODE_ENV=production
  env NODE_ENV=production
  #cd $HOME
  cd /var/www/ilovestage.co.uk
  exec forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk app-api.js
  env PORT=5050
  exec forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk app-api.js
  env PORT=5051
  exec forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk app-api.js
  env PORT=5052
  exec forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log -c "node --harmony" -w --watchDirectory /var/www/ilovestage.co.uk app-api.js
  env PORT=5053
end script
